# delayBoot4.service
#!/bin/sh

[Unit]
Description=Run ARDPI
After=network.target
 
[Service]
ExecStart=/bin/bash /home/pi/pw4/xrunBATMAN4.sh
 
[Install]
WantedBy=multi-user.target
